# etcdlabs-frontend.service
[Unit]
Description=etcdlabs with Docker
Documentation=https://github.com/coreos/etcdlabs

[Service]
Restart=always
RestartSec=5s
TimeoutStartSec=0
LimitNOFILE=40000

ExecStartPre=/usr/bin/docker \
    kill \
    etcdlabs-frontend

ExecStartPre=/usr/bin/docker \
    rm --force \
    etcdlabs-frontend

ExecStartPre=/usr/bin/docker \
    pull \
    quay.io/coreos/etcdlabs:latest

ExecStart=/usr/bin/docker \
    run \
    --net=host \
    --name etcdlabs-frontend \
    --ulimit nofile=262144:262144 \
    quay.io/coreos/etcdlabs:latest \
    yarn \
    start-prod

ExecStop=/usr/bin/docker rm --force etcdlabs-frontend

[Install]
WantedBy=multi-user.target
